# apps_sel_and_conf/cdd/conf/2005-mini-escritorio/icewm/debian/postinst
#! /bin/sh

set -e

case "$1" in
    configure)
        /usr/sbin/update-dpsyco-skel
	#update-rc.d fetchmail-gouda defaults 99
	
	# Configuramos rox-filer

	var=$(cat /etc/passwd | grep :x:1000:1000)
	user=${var/:x:*/}

	echo "<icon x=\"96\" y=\"64\" label=\"Carpeta de Inicio de "$user"\">/home/"$user"</icon>" >> /tmp/.rox_choices/ROX-Filer/pb_mypin
	echo "</pinboard>" >> /tmp/.rox_choices/ROX-Filer/pb_mypin

	# Copiamos los archivos de configuración de rox
	mv /tmp/.rox_choices /home/$user

	# Corregimos los permisos
	chown $user:$user -R /home/$user/.rox_choices

	# Escribimos el script de configuración en el archivo que lanzará
	# rox-filer en modo "desktop" cada vez que se inicie sesión.
	echo "#!/bin/bash" >> /etc/X11/Xsession.d/98rox-filer
	echo "rox-filer -p=/home/\$(whoami)/.rox_choices/ROX-Filer/pb_mypin" >> /etc/X11/Xsession.d/98rox-filer



	# Aloh, ¿Quién está por aquí?
	var=$(cat /etc/passwd | grep :x:1000:1000)
	user=${var/:x:*/}
	
	# Movemos carpeta de documentos a /home/$user
	mv /tmp/documentos /home/$user
	chown $user:$user -R /home/$user/documentos
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)

    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 0
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

