# apps_sel_and_conf/cdd/conf/2005-mini-escritorio/icewm/debian/postinst
#! /bin/sh

set -e

case "$1" in
    configure)
        /usr/sbin/update-dpsyco-skel
	ln -s /usr/X11R6/lib/xdm /etc/X11/xdm
	#update-rc.d fetchmail-gouda defaults 99
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)

    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 0
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

